# bin/conf/sawtooth_identity.service
# Service definition for Hyperledger sawtooth identity tp service
[Unit]
Description=Sawtooth identity tp
Documentation=https://sawtooth.hyperledger.org/docs/1.2/
Wants=network-online.target
After=network-online.target
[Service]
Type=simple
Restart=on-failure
StandardOutput=append:/var/log/sawtooth/sawtooth-identity-sysout.log
StandardError=append:/var/log/sawtooth/sawtooth-identity-syserr.log
ExecStart=identity-tp -vv
[Install]
WantedBy=multi-user.target
