# bin/conf/sawtooth_rest.service
# Service definition for Hyperledger sawtooth rest service
[Unit]
Description=Sawtooth rest server
Documentation=https://sawtooth.hyperledger.org/docs/1.2/
Wants=network-online.target
After=network-online.target
[Service]
Type=simple
Restart=on-failure
StandardOutput=append:/var/log/sawtooth/sawtooth-rest-sysout.log
StandardError=append:/var/log/sawtooth/sawtooth-rest-syserr.log
ExecStart=/usr/bin/sawtooth-rest-api -vvv
[Install]
WantedBy=multi-user.target
