# bin/conf/sawtooth_setting.service
# Service definition for Hyperledger sawtooth settings tp service
[Unit]
Description=Sawtooth setting tp
Documentation=https://sawtooth.hyperledger.org/docs/1.2/
Wants=network-online.target
After=network-online.target
[Service]
Type=simple
Restart=on-failure
StandardOutput=append:/var/log/sawtooth/sawtooth-settings-sysout.log
StandardError=append:/var/log/sawtooth/sawtooth-settings-syserr.log
ExecStart=/usr/bin/settings-tp -vv
[Install]
WantedBy=multi-user.target
