# bin/conf/sawtooth_pbft_engine.service
# Service definition for Hyperledger sawtooth pbft service
[Unit]
Description=Sawtooth pbft
Documentation=https://sawtooth.hyperledger.org/docs/1.2/
Wants=network-online.target
After=network-online.target
[Service]
Type=simple
Restart=on-failure
StandardOutput=append:/var/log/sawtooth/sawtooth-pbft-sysout.log
StandardError=append:/var/log/sawtooth/sawtooth-pbft-syserr.log
ExecStart=pbft-engine -C tcp://127.0.0.1:5050
[Install]
WantedBy=multi-user.target